# src/cron.d/default
*/1 * * * * root /usr/local/bin/php /www/e-chat/yii daemon/logout
